# src/res/hypnos.service.in
[Unit]
Description=Hypnos idle daemon service
Documentation=https://github.com/axos-project/hypnos
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=@BIN_PATH@ -c @HOME@/.sleex/hypnos.json
Restart=on-failure

[Install]
WantedBy=default.target